# lab3/crontab
*/5 * * * 6 /root/lab3/i
*/1 * * * 3 /root/lab3/i
*/1 * * * thu /root/lab3/i
*/1 crontab i ii iii ~report report test crontab i ii iii ~report report test crontab i ii iii ~report report test thu /root/lab3/i
*/2 crontab i ii iii ~report report test crontab i ii iii ~report report test crontab i ii iii ~report report test thu /root/lab3/i
*/2 crontab i ii iii ~report report test crontab i ii iii ~report report test crontab i ii iii ~report report test wed /root/lab3/i
*/2 crontab i ii iii ~report report test crontab i ii iii ~report report test crontab i ii iii ~report report test wed /root/lab3/i
* * * * wed ~/lab3/i
* * * * sut ~/lab3/i
